# infra/feather.supervisor.conf
; Feather — supervisord config
; Copy to /etc/supervisor/conf.d/ and reload:
;   sudo cp feather.supervisor.conf /etc/supervisor/conf.d/feather.conf
;   supervisorctl reread && supervisorctl update
;
; Commands:
;   supervisorctl status feather
;   supervisorctl restart feather
;   supervisorctl tail -f feather

[program:feather]
command=%(here)s/run.sh
directory=%(here)s
environment=PORT="4870",HOME="/home/user"
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/tmp/feather.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_redirect=true
stopwaitsecs=10
